# Server/mobileterminal.service
# Mobile Terminal Service with Session Persistence
#
# Install location: /etc/systemd/system/mobileterminal.service
#
# After installing:
#   sudo systemctl daemon-reload
#   sudo systemctl restart mobileterminal
#   sudo systemctl enable mobileterminal
#

[Unit]
Description=Mobile Terminal (ttyd with tmux session persistence)
After=network.target

[Service]
Type=simple
User=root
Environment=TERM=xterm-256color

# ttyd with tmux wrapper for session persistence
# -W: writable (allows input)
# -p 7681: port number
# -i 127.0.0.1: bind to localhost only (use reverse proxy for external access)
# -t fontSize=24: default font size for mobile
# -t titleFixed=Claude Mobile: window title
ExecStart=/usr/bin/ttyd -W -p 7681 -i 127.0.0.1 -t fontSize=24 -t titleFixed="Claude Mobile" /usr/local/bin/ttyd-tmux.sh

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
